# installSettings/startMailScript.service
[Unit]
Description=Start django server
[Service]
ExecStart=/etc/startMailScript.sh
[Install]
WantedBy=multi-user.target
