# installSettings/startDjango.service
[Unit]
Description=Start django server
[Service]
ExecStart=/etc/autostart.sh
[Install]
WantedBy=multi-user.target
